# jobs/choco.service
[Unit]
Description=choco - CHORD Config Orchestrator
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
WorkingDirectory=/opt/choco
ExecStart=/opt/choco/.venv/bin/choco /etc/choco/config.yaml
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
