# jobs/choco-eop-broadcast.timer
[Unit]
Description=EOP table broadcast (3x daily)

[Timer]
OnCalendar=*-*-* 06,14,22:00:00
Persistent=true

[Install]
WantedBy=timers.target
